# templates/99-ncdata.rules
# homecloud: auto-mount SSD and restart services when the external drive reappears.
# Installed to /etc/udev/rules.d/99-ncdata.rules
#
# The LABEL= check matches the SSD label set during install (default: ncdata).
# When the device is added (plugged in), we mount it and restart Docker + Nextcloud.

ACTION=="add", ENV{ID_FS_LABEL}=="{{SSD_LABEL}}", RUN+="/bin/systemctl start ncdata-replug.service"
