# templates/homecloud-bot.service
[Unit]
Description=Home Cloud Telegram Bot (Immich)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt
ExecStart=/opt/homecloud-bot-env/bin/python3 /opt/homecloud-bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
